# systemd/mintmark@.service
[Unit]
Description=Spool Markdown to an Epson TM-U220B receipt printer
RefuseManualStart=true

[Service]
Type=simple
Environment=MINTMARK_DEVICE=/dev/usb/lp0
ExecStart=mintmark --lock-file ${RUNTIME_DIRECTORY}/lock ${MINTMARK_DEVICE}
User=mintmark
Group=lp
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
RuntimeDirectory=mintmark
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=true
PrivateTmp=true
PrivateNetwork=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
StandardInput=socket
StandardOutput=null
StandardError=socket
